# systemd/user/ollama.service
[Unit]
Description= Ollama Service

[Service]
Type=oneshot
ExecStart=/usr/bin/ollama serve

[Install]
WantedBy=default.target
