# systemd/user/update-wallpaper.timer
[Unit]
Description=Update user wallpaper onn hourly basis

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
